# strategy-lab.service
[Unit]
Description=Trading Guru Strategy Lab v1.0 — Paper Testing + Skill Training + Auto-Switch
After=network.target trading-engine.service
Wants=trading-engine.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading_engine
ExecStart=/usr/bin/python3 /home/ubuntu/trading_engine/strategy_lab.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/trading_engine/strategy_lab.log
StandardError=append:/home/ubuntu/trading_engine/strategy_lab.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
