# supervisor.conf
[program:l99bot]
directory=/home/YOUR_USER/ai-trading-championship/l99_bot
command=/home/YOUR_USER/ai-trading-championship/l99_bot/venv/bin/python live_bot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/l99bot.err.log
stdout_logfile=/var/log/l99bot.out.log
user=YOUR_USER
environment=TESTNET="true",LIVE_TRADING="false"

; Deploy:
;   sudo cp supervisor.conf /etc/supervisor/conf.d/l99bot.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start l99bot
