# canary/champion-battle.service
[Unit]
Description=Champion Battle — CMO Gods Level (TITAN/VELOCITY/SENTINEL) [DigitalOcean]
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/canary
Environment=CANARY_ROOT=/root/canary
ExecStart=/usr/bin/python3 /root/canary/champion_battle.py
Restart=always
RestartSec=5
StandardOutput=append:/root/canary/runtime/champion_battle.systemd.log
StandardError=append:/root/canary/runtime/champion_battle.systemd.log

[Install]
WantedBy=multi-user.target
