# canary/coach.service
[Unit]
Description=Trading Guru Team Coach / Coordinator (smart dynamic rotation) [DigitalOcean]
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/canary
Environment=CANARY_BASE=/root
ExecStart=/usr/bin/python3 /root/canary/coach.py
Restart=always
RestartSec=5
StandardOutput=append:/root/canary/runtime/coach.systemd.log
StandardError=append:/root/canary/runtime/coach.systemd.log

[Install]
WantedBy=multi-user.target
