# paper_battle/tradingguru-arena.service
[Unit]
Description=⚔️ TradingGuru Shadow Arena — feature enricher (weather + boss + crowd)
Documentation=file:///root/tradingguru-empire/governance/SHADOW_ARENA_FEATURES_v1.md
Documentation=file:///root/tradingguru-empire/arena_shadow_runner.py
After=network-online.target paper-arena.service
Wants=paper-arena.service

# ── ISOLATION GUARDS ────────────────────────────────────────────────────
# This service is paper-only + visualization-only. It writes arena_features.json.
# Refuses to start if Layer 1 mutation detected (paranoid):
ConditionPathExists=/root/tradingguru-empire/canary/canary_strategy.py
ConditionPathExists=/root/.l99/protection_halt.json

[Service]
Type=simple
# Hard safety env (refuses to start without these):
Environment="TG_DISABLE_REAL_ORDERS=1"
Environment="TG_PAPER_EXECUTION=1"
Environment="LIVE_ORDERS=0"
# Optional: publish to nginx-served path so frontend can fetch /api/battle/arena-features.json
Environment="ARENA_FEATURES_PUBLISH_PATH=/var/www/ai-trading-championship/dist/api/battle/arena-features.json"

ExecStart=/usr/bin/python3 /root/tradingguru-empire/arena_shadow_runner.py \
    --mode shadow --market real --execution paper \
    --agents HUNTER,RISK,ALPHA,REGIME,EXECUTOR,RECOVERY,SKEPTIC,CHAMPION \
    --enable-xp --enable-commentary --enable-rivalries \
    --enable-weather --enable-boss-events --enable-crowd \
    --feed terminal.json \
    --scan-interval 10

WorkingDirectory=/root/tradingguru-empire
Restart=on-failure
RestartSec=30
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
