# systemd/d6-autotrigger.service
[Unit]
Description=D6 Phase A binding-run auto-trigger (one-shot)
Documentation=https://github.com/RolandGasparyan/ai-trading-championship
After=microstructure-collector.service
Wants=microstructure-collector.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/ai-trading-championship
EnvironmentFile=-/etc/default/trading-bot-telegram
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/ai-trading-championship/venv/bin/python \
          /var/www/ai-trading-championship/d6_autotrigger.py
StandardOutput=append:/var/log/d6_autotrigger.log
StandardError=append:/var/log/d6_autotrigger.log
TimeoutStartSec=2100
MemoryMax=512M

[Install]
WantedBy=multi-user.target
