# systemd/l2-depth-collector.service
[Unit]
Description=L2 Depth Collector — Gate.io spot orderbook ladder capture (D3 research)
Documentation=https://github.com/RolandGasparyan/ai-trading-championship
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trading-bot
Group=trading-bot
WorkingDirectory=/opt/ai-trading-championship
ExecStart=/usr/bin/python3 /opt/ai-trading-championship/l2_depth_collector.py \
    --pairs BTC_USDT ETH_USDT SOL_USDT XRP_USDT AVAX_USDT \
    --out-dir /var/log/l2_depth \
    --snap-hz 1.0
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
EnvironmentFile=-/etc/default/l2-depth-collector

# Resource caps (lightweight — single-process asyncio)
MemoryHigh=512M
MemoryMax=1G
CPUQuota=50%

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/l2_depth

[Install]
WantedBy=multi-user.target
