# systemd/production-monitor.service
[Unit]
Description=Production health monitor (one-shot, fired by timer)
Documentation=https://github.com/RolandGasparyan/ai-trading-championship
After=bot-control.service trading-bot@.service microstructure-collector.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/ai-trading-championship
EnvironmentFile=-/etc/default/trading-bot-telegram
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/ai-trading-championship/venv/bin/python \
          /var/www/ai-trading-championship/production_monitor.py
StandardOutput=append:/var/log/production_monitor.log
StandardError=append:/var/log/production_monitor.log
TimeoutStartSec=120
MemoryMax=64M

[Install]
WantedBy=multi-user.target
