# production_engine/config/trading-engine.service
[Unit]
Description=AI Trading Engine - UNCRISHABLE
After=network.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/trading-engine
ExecStart=/usr/bin/python3 /opt/trading-engine/engine/main.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/trading-engine/logs/engine.log
StandardError=append:/opt/trading-engine/logs/errors.log

# Environment
EnvironmentFile=/opt/trading-engine/.env

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trading-engine

[Install]
WantedBy=multi-user.target
